# kiosk.service
[Unit]
Description=Kiosk Django Application (Local - No Internet)
After=network.target

[Service]
Type=notify
# برای لوکال می‌توانی user خودت را بگذاری یا root
User=root
Group=root
WorkingDirectory=/opt/kiosk
Environment="PATH=/opt/kiosk/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings"
# بدون نیاز به network-online.target برای لوکال
ExecStart=/opt/kiosk/venv/bin/gunicorn \
    --bind 127.0.0.1:8000 \
    --workers 4 \
    --timeout 300 \
    config.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
# Timeout بیشتر برای لوکال
TimeoutStartSec=300
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target

